# hydrocore3.1/packaging/hydrocore-watchdog.service
[Unit]
Description=HydroCore HTTP health check
After=hydrocore.service

[Service]
Type=oneshot
ExecStart=/usr/bin/hydrocore-watchdog
