# hydrocore3.1/packaging/hydrocore-watchdog.timer
[Unit]
Description=Run HydroCore health check every minute

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s
AccuracySec=10s
Unit=hydrocore-watchdog.service

[Install]
WantedBy=timers.target
