# hydrocore3.1/packaging/hydrocore-reboot.service
[Unit]
Description=Reboot HydroCore device from maintenance UI

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'rm -f /var/lib/hydrocore/runtime/reboot.request; sleep 1; /usr/bin/systemctl reboot'
